# src/server/supervisord.conf
# src/server/supervisord.conf

[supervisord]
nodaemon=true
user=root
stdout_logfile=/var/log/supervisor/supervisord.log
stderr_logfile=/var/log/supervisor/supervisord_err.log

[program:main-server]
command=uvicorn main.app:app --host 0.0.0.0 --port 5000 --workers 2
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/main-server.log
stderr_logfile=/var/log/supervisor/main-server_err.log

[program:celery-worker]
command=celery -A server.workers.celery_app worker --loglevel=info --pool=solo
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/celery-worker.log
stderr_logfile=/var/log/supervisor/celery-worker_err.log

[program:celery-beat]
command=celery -A server.workers.celery_app beat --loglevel=info
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/celery-beat.log
stderr_logfile=/var/log/supervisor/celery-beat_err.log

# --- MCP Servers ---

[program:mcp-accuweather]
command=python -m mcp-hub.accuweather.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-accuweather.log
stderr_logfile=/var/log/supervisor/mcp-accuweather_err.log

[program:mcp-chat-tools]
command=python -m mcp-hub.chat_tools.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-chat-tools.log
stderr_logfile=/var/log/supervisor/mcp-chat-tools_err.log

[program:mcp-gcal]
command=python -m mcp-hub.gcal.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-gcal.log
stderr_logfile=/var/log/supervisor/mcp-gcal_err.log

[program:mcp-gdocs]
command=python -m mcp-hub.gdocs.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-gdocs.log
stderr_logfile=/var/log/supervisor/mcp-gdocs_err.log

[program:mcp-gdrive]
command=python -m mcp-hub.gdrive.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-gdrive.log
stderr_logfile=/var/log/supervisor/mcp-gdrive_err.log

[program:mcp-github]
command=python -m mcp-hub.github.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-github.log
stderr_logfile=/var/log/supervisor/mcp-github_err.log

[program:mcp-gmail]
command=python -m mcp-hub.gmail.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-gmail.log
stderr_logfile=/var/log/supervisor/mcp-gmail_err.log

[program:mcp-gmaps]
command=python -m mcp-hub.gmaps.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-gmaps.log
stderr_logfile=/var/log/supervisor/mcp-gmaps_err.log

[program:mcp-google-search]
command=python -m mcp-hub.google_search.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-google-search.log
stderr_logfile=/var/log/supervisor/mcp-google-search_err.log

[program:mcp-gshopping]
command=python -m mcp-hub.gshopping.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-gshopping.log
stderr_logfile=/var/log/supervisor/mcp-gshopping_err.log

[program:mcp-gsheets]
command=python -m mcp-hub.gsheets.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-gsheets.log
stderr_logfile=/var/log/supervisor/mcp-gsheets_err.log

[program:mcp-gslides]
command=python -m mcp-hub.gslides.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-gslides.log
stderr_logfile=/var/log/supervisor/mcp-gslides_err.log

[program:mcp-journal]
command=python -m mcp-hub.journal.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-journal.log
stderr_logfile=/var/log/supervisor/mcp-journal_err.log

[program:mcp-news]
command=python -m mcp-hub.news.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-news.log
stderr_logfile=/var/log/supervisor/mcp-news_err.log

[program:mcp-notion]
command=python -m mcp-hub.notion.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-notion.log
stderr_logfile=/var/log/supervisor/mcp-notion_err.log

[program:mcp-progress-updater]
command=python -m mcp-hub.progress_updater.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-progress-updater.log
stderr_logfile=/var/log/supervisor/mcp-progress-updater_err.log

[program:mcp-quickchart]
command=python -m mcp-hub.quickchart.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-quickchart.log
stderr_logfile=/var/log/supervisor/mcp-quickchart_err.log

[program:mcp-slack]
command=python -m mcp-hub.slack.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp-slack.log
stderr_logfile=/var/log/supervisor/mcp-slack_err.log

[program:healthcheck]
command=python ping_waha.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/healthcheck.log
stderr_logfile=/var/log/supervisor/healthcheck_err.log